# cvsnt/cvsnt-2.5.05.3744/redhat/cvslockd.service
[Unit]
Description=Handles file level locking for the cvsnt server
After=network.target

[Service]
Type=simple
LOCKFILE=/var/lock/subsys/cvslockd
PIDFile=/run/cvslockd.pid
OOMScoreAdjust=-100
EnvironmentFile=-/etc/sysconfig/cvslockd
ExecStart=/usr/bin/cvslockd -n
KillMode=control-group
Restart=on-failure

[Install]
WantedBy=multi-user.target

